# supervisor/conf.d/rundeck.conf
[program:rundeck]
command=/etc/init.d/rundeckd start
